# userpatches/overlay/dietpi-auto-install.service
[Unit]
Description=DietPi Auto-Install (First Boot)
After=local-fs.target network-online.target
Wants=network-online.target
Before=multi-user.target
ConditionPathExists=!/boot/dietpi/.auto-install-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dietpi-auto-install.sh
RemainAfterExit=yes
StandardOutput=journal+console
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
